# etc/cron.d/shadowsocks
@reboot  kk /opt/shadowsocks/shadowsocks-Linux64bit
